# scripts/systemd/lotus-worker-p2s.target
[Unit]
Description=Lotus Worker P2S
Requires=lotus-worker-p2@0.service lotus-worker-p2@1.service

[Install]
WantedBy=multi-user.target
