# scripts/systemd/lotus-worker-p1s.target
[Unit]
Description=Lotus Worker P1S
Requires=lotus-worker-p1@0.service

[Install]
WantedBy=multi-user.target
