# scripts/systemd/logrotate/lotus.cron
5 0 * * *  /usr/sbin/logrotate /etc/logrotate.d/lotus