# scripts/lotus-daemon.service
[Unit]
Description=Lotus Daemon
After=network-online.target
Requires=network-online.target

[Service]
User=lotus
Group=lotus
ExecStartPre=/usr/bin/lotus-init.sh
ExecStart=/usr/bin/lotus daemon
Environment=LOTUS_PATH=/var/lib/lotus
Restart=always
RestartSec=30
TimeoutSec=infinity

[Install]
WantedBy=multi-user.target
